# radeon/RadeonSwitchOff.service
[Unit]
Description=Switch off Radeon discrete card

[Service]
ExecStart=/usr/local/bin/switchoff

[Install]
WantedBy=multi-user.target
